# auto-internet-speed-checker/auto-internet-speed-checker.service
[Service]
ExecStart=/usr/local/bin/systemd_scheduler.sh